# src/ssh_tunnel.init
#!/bin/sh /etc/rc.common

START=99
STOP=10

PROCD_DEBUG=1
INIT_TRACE=1 

PROG="/root/ssh_tunnel.sh"

start_service() {
    procd_open_instance
    procd_set_param command /bin/sh /root/ssh_tunnel.sh
    procd_close_instance
}

stop_service() {
    "$PROG" stop
}

restart() {
    stop
    sleep 2
    start
}
