# src/ssh_tunnel.init
#!/bin/sh /etc/rc.common

START=99
STOP=10

PROCD_DEBUG=1
INIT_TRACE=1 /etc/init.d/ssh_tunnel start

PROG=/root/ssh_tunnel.sh

start_service() {
    procd_set_param command "$PROG" start
}

stop_service() {
    "$PROG" stop
}

restart() {
    stop
    sleep 2
    start
}
